# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=docker-php-entrypoint php-fpm
autostart=true
autorestart=true

[program:queue-worker]
command=php artisan queue:work --queue=default --sleep=3 --tries=3
autostart=true
autorestart=true
